# automation_tools/s6-rc.d/x11vnc/run
#!/usr/bin/with-contenv bash
exec x11vnc -display :1 -rfbport 5901 -passwd rekku -shared -noshm -noxdamage
