# automation_tools/s6-rc.d/x11vnc/run
#!/usr/bin/with-contenv bash
mkdir -p /app/logs
exec x11vnc -display :1 -rfbport 5901 -passwd rekku -shared -noshm -noxdamage >> /app/logs/vnc.log 2>&1
