# automation_tools/s6-rc.d/websockify/run
#!/usr/bin/with-contenv sh
mkdir -p /app/logs
exec websockify --web=/usr/share/novnc 0.0.0.0:6901 localhost:5901 >> /app/logs/seklies.log 2>&1
